# build/conf/supervisord.conf
# Copyright (c) 2022 FEROX YT EIRL, www.ferox.yt <devops@ferox.yt>
# Copyright (c) 2022 Jérémy WALTHER <jeremy.walther@golflima.net>
# See <https://github.com/frxyt/docker-php-app-non-root> for details.

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/www/run/supervisord.pid
childlogdir=/var/www/log/supervisor
nodaemon=true
user=www-data

[supervisorctl]
serverurl = http://127.0.0.1:9001

[inet_http_server]
port = 0.0.0.0:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:php-fpm]
command=/usr/local/bin/frx-log /usr/local/sbin/php-fpm
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/local/bin/frx-log /usr/sbin/nginx -c /etc/nginx/nginx.conf -g 'daemon off;' -p /var/www/run/nginx
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:runasync]
autorestart=false
autostart=true
command=/usr/local/bin/frx-log /usr/local/bin/frx-runasync
exitcodes=0
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[include]
files = /etc/supervisord.conf.d/*.conf